# vigilante.service
[Unit]
Description=Test downloads

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/pi/vigilante/vigilante.py
